# infra/doc/jail/data/nslcd/rc.local
chroot /compat/jammy /usr/sbin/service nslcd start
